# init/pcmad.service.in
[Unit]
Description=Page Cache My Assets daemon

[Service]
ExecStart=@bindir@/pcmad
OOMScoreAdjust=-1000
StandardOutput=syslog
UMask=000

[Install]
WantedBy=multi-user.target
